# ansible-collections/grnet/jitsi/roles/jitsi_meet/templates/jitsi-prometheus-exporter.service
[Unit]
Description=Jitsi Prometheus Exporter
After=network.target

[Service]
ExecStart=/usr/local/bin/jitsi-prometheus-exporter {{ jitsi_prometheus_exporter_port }}
Restart=on-failure

[Install]
WantedBy=multi-user.target
